# src/deb/control/postinst
#!/bin/bash
set -e
# 开始启动
/etc/init.d/zrlogwww start
# 开机启动
update-rc.d zrlogwww defaults